# batman-gwsel/crontab
# update Gatewayselection
*/5 * * * * 	root /usr/local/bin/dyn_announce_gw_bw.sh

